# net-p2p/btsync/files/init.d/btsync
#!/sbin/runscript
# Copyright (C) 2013 Jonathan Vasquez <jvasquez1011@gmail.com>
# Distributed under the terms of the Simplified BSD License.

SYNC_NAME="btsync"
SYNC_PATH="/opt/${SYNC_NAME}/"
SYNC_BINARY="${SYNC_PATH}/${SYNC_NAME}"
SYNC_PIDFILE="/var/run/${SYNC_NAME}.pid"

depend() {
    need localmount net
    after bootmisc
}

start() {
    if [ ! -d /run/btsync ]; then
        mkdir /run/btsync
        chown "${BTSYNC_USER}:${BTSYNC_GROUP}" /run/btsync
    fi

    ebegin "Starting ${SYNC_NAME}"
    start-stop-daemon \
        --start \
        --pidfile "${PIDFILE}" \
        --user "${BTSYNC_USER}" \
        --group "${BTSYNC_GROUP}" \
        --exec "${SYNC_BINARY}" \
        -- --config "/etc/btsync.conf"
    eend $?
}

stop() {
	ebegin "Stopping ${SYNC_NAME}"

	start-stop-daemon --stop --exec "${SYNC_BINARY}" \
	    --pidfile "${SYNC_PIDFILE}"

	eend $?
}

reload() {
	ebegin "Reloading ${SYNC_NAME}"

	start-stop-daemon --signal HUP --exec "${SYNC_BINARY}" \
	    --pidfile "${SYNC_PIDFILE}"

	eend $?
}
