# net-p2p/btsync/files/btsync.service
[Service]
Type=forking
PIDFile=/opt/btsync/.sync
ExecStart=/opt/btsync/btsync

[Install]
WantedBy=multi-user.target

