# net-p2p/btsync/files/btsync.service
[Service]
Type=simple
ExecStart=/opt/btsync/btsync --nodaemon

[Install]
WantedBy=multi-user.target

