# net-p2p/btsync/files/btsync.service
[Service]
Type=dbus
BusName=org.btsync.daemon
ExecStart=/opt/btsync

[Install]
WantedBy=multi-user.target

